# deploy/go-carbon.service
[Unit]
Description=Golang implementation of Graphite/Carbon server.
Documentation=https://github.com/lomik/go-carbon
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/go-carbon -config /etc/go-carbon/go-carbon.conf -pidfile /var/run/go-carbon.pid -daemon
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=USR2
Restart=on-failure

TimeoutStopSec=600
LimitNOFILE=55555
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target